# config/systemd/vpn-namespace.service
[Unit]
Description=VPN Network Namespace for qBittorrent
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/home/anon/nas-media-server/scripts/vpn-namespace-setup.sh setup
ExecStop=/home/anon/nas-media-server/scripts/vpn-namespace-setup.sh cleanup

[Install]
WantedBy=multi-user.target
