# configs/livetv-update.service
[Unit]
Description=Update IPTV channels and EPG data
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=anon
ExecStart=/home/anon/nas-media-server/scripts/update-iptv.sh
ExecStartPost=/home/anon/nas-media-server/scripts/update-epg.sh
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target
